# systemd/wopr-production-deployer.service
[Unit]
Description=WOPR Production Deployer (mesh-native)
After=network.target

[Service]
User=root
ExecStart=/usr/bin/sleep infinity
Restart=always

[Install]
WantedBy=multi-user.target
